# container_files/etc/supervisord.d/haproxy.ini
[program:haproxy]
command = haproxy -f /etc/haproxy/haproxy.cfg
stderr_logfile = syslog
stdout_logfile = syslog